# infra/systemd/voxaflow-ngrok.service
[Unit]
Description=VoxaFlow ngrok tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/voxaflow/ngrok.env
ExecStart=/usr/local/bin/ngrok http ${NGROK_LOCAL_PORT} --log stdout --log-format json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
